# systemd/agent-linux.service
[Unit]
Description=MYZEL Agent Linux (Lite)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# IMPORTANT:
# Edit these two lines after copying to the device, if needed.
User=pi
WorkingDirectory=/home/pi/agent-linux

Environment=PYTHONUNBUFFERED=1

# Use .venv consistently (matches our Linux install steps)
ExecStart=/home/pi/agent-linux/.venv/bin/python /home/pi/agent-linux/app.py

Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
