# systemd/agent-linux.service
[Unit]
Description=MYZEL Agent Linux (Lite)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
Restart=always
RestartSec=2

# ---- EDIT THESE 3 LINES ON THE TARGET LINUX MACHINE ----
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/agent-linux
ExecStart=/home/YOUR_USER/agent-linux/.venv/bin/python /home/YOUR_USER/agent-linux/app.py
# --------------------------------------------------------

[Install]
WantedBy=multi-user.target
